# domeneshop/download.crontab.template
00 22 * * *  bash -c '<PATH>; #<MARKER>' 
